# root/etc/systemd/system/dnscrypt-proxy.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dnscrypt-proxy --config /etc/dnscrypt-proxy/custom/dnscrypt-proxy.toml
