# .config/systemd/user/gpg-agent.service.d/override.conf
[Install]
WantedBy=default.target

# vim: set ft=systemd :
