# root/etc/systemd/system/usbguard.service.d/override.conf
[Service]
# Workaround for https://github.com/USBGuard/usbguard/issues/261
Restart=always
RestartSec=5
