# root/etc/systemd/system/updatedb.service.d/override.conf
# Prevent updatedb from indexing files in snapshots
# https://wiki.archlinux.org/title/Snapper#updatedb
[Service]
ExecStart=
ExecStart=/usr/bin/updatedb --add-prunenames '.snapshots'
